# infra/systemd/athens-keeper.service
[Unit]
Description=Athens EHS System Keeper
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/athens
ExecStart=/var/www/athens/athens-keeper.sh start
ExecStop=/var/www/athens/athens-keeper.sh stop
Restart=always
RestartSec=10
Environment=ATHENS_BACKEND_PORT=8001
Environment=VITE_PORT=3000

[Install]
WantedBy=multi-user.target