# infra/systemd/rayzen-uvicorn.service
[Unit]
Description=Rayzen ASGI (uvicorn)
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/Rayzen/athenas-backend
EnvironmentFile=/var/www/Rayzen/athenas-backend/.env
Environment=PYTHONUNBUFFERED=1
Environment=DJANGO_SETTINGS_MODULE=config.settings
ExecStart=/var/www/Rayzen/athenas-backend/venv/bin/uvicorn config.asgi:application --host 127.0.0.1 --port 8002 --workers 1 --proxy-headers --forwarded-allow-ips=*
Restart=always
RestartSec=5
UMask=0022

[Install]
WantedBy=multi-user.target
